# docs/launchd/com.ecc.hygiene.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.ecc.hygiene.plist — ECC memory hygiene worker.

  Runs scripts/hygiene-worker.py once per hour. The plist uses
  launchd's StartInterval (3600 s) to schedule the work; the
  worker is invoked with --once so it exits after a single pass.
  launchd respawns it on the schedule.

  The worker invokes memory.run_hygiene_pass() (defined in
  init_schema.sql) to clean up retrieval logs, expired
  typed_memory rows, orphan memory_edges, the audit log,
  stale pending approvals, and old dream runs. It also records
  per-pass metrics to the local telemetry database.

  Replace <PYTHON_PATH>, <REPO_ROOT>, and <HOME> below before loading.
  The DATABASE_URL must authenticate as a trusted database service
  identity such as agent_memory_service; setting app.current_role alone
  is not sufficient under forced RLS.

  To install:
    cp com.ecc.hygiene.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.ecc.hygiene.plist

  To uninstall:
    launchctl unload ~/Library/LaunchAgents/com.ecc.hygiene.plist
    rm ~/Library/LaunchAgents/com.ecc.hygiene.plist

  Logs:
    stdout: ~/Library/Logs/com.ecc.hygiene.log
    stderr: ~/Library/Logs/com.ecc.hygiene.err
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ecc.hygiene</string>
    <key>ProgramArguments</key>
    <array>
        <string>&lt;PYTHON_PATH&gt;</string>
        <string>&lt;REPO_ROOT&gt;/scripts/hygiene-worker.py</string>
        <string>--once</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DATABASE_URL</key>
        <string>postgresql://agent_memory_service@/agent_memory</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>&lt;HOME&gt;/Library/Logs/com.ecc.hygiene.log</string>
    <key>StandardErrorPath</key>
    <string>&lt;HOME&gt;/Library/Logs/com.ecc.hygiene.err</string>
</dict>
</plist>
