# docs/launchd/com.agent_architecture.dream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.agent_architecture.dream.plist — biweekly batch memory curation.

  Runs the memory-dream skill (skills/memory-dream/scripts/dream.py)
  on the 1st and 15th of each month at 11:00 local time. Produces a
  staging set in memory.dream_proposals; does NOT auto-apply.

  Workflow after this runs:
    1. Review the staging via `notify_review.py` (Telegram digest)
       or directly:   python skills/memory-dream/scripts/dream.py runs
    2. Inspect:       python skills/memory-dream/scripts/dream.py diff <run_id>
    3. Adopt:         python skills/memory-dream/scripts/dream.py adopt <run_id>
       or discard:    python skills/memory-dream/scripts/dream.py discard <run_id>

  Cadence: `0 11 1,15 * *` (1st and 15th at 11:00). Same pattern
  as hermes-dream's weekly cron.

  Replace <PYTHON_PATH>, <HOME>, and <USER> below before loading.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent_architecture.dream</string>
    <key>ProgramArguments</key>
    <array>
        <string>&lt;PYTHON_PATH&gt;</string>
        <string>&lt;HOME&gt;/agent_architecture/skills/memory-dream/scripts/dream.py</string>
        <string>run</string>
        <string>--user-id</string>
        <string>owner:&lt;USER&gt;</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>11</integer>
        <key>Day</key>
        <array>
            <integer>1</integer>
            <integer>15</integer>
        </array>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DATABASE_URL</key>
        <string>postgresql:///agent_memory</string>
        <key>ACTOR_ID</key>
        <string>owner:&lt;USER&gt;</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>StandardOutPath</key>
    <string>&lt;HOME&gt;/Library/Logs/com.agent_architecture.dream.log</string>
    <key>StandardErrorPath</key>
    <string>&lt;HOME&gt;/Library/Logs/com.agent_architecture.dream.err</string>
</dict>
</plist>
